# deploy/com.gncausreport.bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- macOS launchd: Socket Mode 봇 상시 데몬(승인/반려·이메일 발송·주간주제 버튼 처리).
     설치: cp deploy/com.gncausreport.bot.plist ~/Library/LaunchAgents/ &&
           launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.gncausreport.bot.plist
     (Linux 서버 배포는 deploy/gncausreport-bot.service systemd 유닛 사용) -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gncausreport.bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/limsungmin/gncausreport/.venv/bin/python</string>
        <string>/Users/limsungmin/gncausreport/bot/server.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/limsungmin/gncausreport</string>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/limsungmin/gncausreport/logs/bot.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/limsungmin/gncausreport/logs/bot.err.log</string>
</dict>
</plist>
